# service/packages/debian/xroad-fileservice.upstart
# Place in /etc/init

description "File serice on X-Road"
author "Sami Kallio <sami.kallio@gofore.com>"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5
umask 007
setuid xroad-fileservice
setgid xroad-fileservice
console log

script

exec java -Dserver.port=8070 -jar /usr/lib/xroad-fileservice/xroad-fileservice.jar

end script
