# service/packages/redhat/SOURCES/xroad-fileservice.service
[Unit]
Description=X-Road file service
After=network.target

[Service]
User=xroad-fileservice
Group=xroad-fileservice
ExecStart=/usr/bin/java -jar /usr/lib/xroad-fileservice/xroad-fileservice.jar
Restart=on-failure

[Install]
WantedBy=multi-user.target

